# recipes-extended/retroarch-autostart-service/files/retroarch.service.in
[Unit]
Description=RetroArch

[Service]
ExecStart=/usr/bin/retroarch
Restart=on-failure
User=@@RETRO_USER_NAME@@

[Install]
WantedBy=multi-user.target
